# deploy/host/nova-runtime-control.service
[Unit]
Description=Nova Runtime Control
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=nova
Group=nova
WorkingDirectory=/home/nova/studio-nova/apps/nova-runtime-control
EnvironmentFile=-/home/nova/studio-nova/apps/nova-runtime-control/.env.local
Environment=KUBECONFIG=/home/nova/.kube/config
Environment=NODE_ENV=production
ExecStart=/home/nova/.vite-plus/bin/node --experimental-strip-types src/index.ts
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=read-only
ProtectSystem=full
ReadWritePaths=/home/nova/.kube

[Install]
WantedBy=multi-user.target
